# init.d/openstack_keystone
#!/sbin/runscript
# Copyright 1999-2015 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

pidfile=/var/run/keystone.pid
command=$OPENSTACK_PREFIX/bin/keystone-all
required_files="${KEYSTONE_CONF:-/etc/keystone/keystone.conf}"
start_stop_daemon_args="--background --make-pidfile" # --user ${KEYSTONE_USER:-keystone}"
command_args="--log-file=$OPENSTACK_PREFIX/log/keystone.log"

depend() {
        need mysql
}
