# init.d/openstack_nova
#!/sbin/runscript
# Copyright 1999-2015 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

pidfile=/var/run/nova.pid
command=$OPENSTACK_PREFIX/bin/nova-all
required_files="${KEYSTONE_CONF:-/etc/nova/nova.conf}"
start_stop_daemon_args="--background --make-pidfile" # --user ${KEYSTONE_USER:-keystone}"

depend() {
        need mysql
}
