# init.d/openstack_neutron-server
#!/sbin/runscript
# Copyright 1999-2015 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

pidfile=/var/run/neutron-server.pid
command="$OPENSTACK_PREFIX/bin/neutron-server"
command_args="--config-file /etc/neutron/neutron.conf --config-file /etc/neutron/plugins/ml2/ml2_conf.ini"
required_files="${NEUTRON_CONF:-/etc/neutron/neutron.conf}"
start_stop_daemon_args="--background --make-pidfile --user ${NEUTRON_USER:-root}"

depend() {
        need mysql
}
